# dot_config/systemd/user/aur-update-notify.service
[Unit]
Description=Check repo and AUR updates and notify

[Service]
Type=oneshot
ExecStart=%h/.local/bin/aur-update-notify
