# dot_config/systemd/user/kdeconnect-keepalive.service
[Unit]
Description=KDE Connect Keepalive Refresh
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/kdeconnect-cli --refresh
